# cluster-management/morello-worker@.service
[Unit]
Description=Morello work-queue and board-runner
After=dev-morello-%i-tty\x2dmcc.device
After=dev-morello-%i-tty\x2dap0.device
BindsTo=dev-morello-%i-tty\x2dmcc.device
BindsTo=dev-morello-%i-tty\x2dap0.device

[Install]
WantedBy=multi-user.target

[Service]
Type=exec
User=morello-auto
Group=morello-auto

# The executor is one-shot; restart it by default.
Restart=always

# The executor exits 42 when it receives an explicit shutdown request.  Those
# should stop it from reentering the pool (RestartPreventExitStatus) and should
# not count as failure (SuccessExitStatus).
RestartPreventExitStatus=42
SuccessExitStatus=42

RuntimeDirectory=morello/%I
Environment=MORELLO_HOSTNAME=%I

# Set $MORELLO_SCRIPTS with a drop-in [Service] Environment=
ExecStart=/bin/sh ${MORELLO_SCRIPTS}/work-bus/executor/wrapper.sh
